# config/smarthome_node.conf
;/etc/supervisor/conf.d/smarthome.conf
[program:smarthome_node]
user=smarthomenode
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
