# config/smarthome_node.ini
[program:smarthome_node]
user=smarthomenode
directory=/path/to/your/app
command=/path/to/your/app/bin/boot.sh 
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d